# onedice-bot.service
[Unit]
Description=OneDice Telegram Bot Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/onedice
ExecStart=/usr/bin/python3 main-simple.py
Restart=always
RestartSec=10
Environment=BOT_TOKEN=8360132125:AAEWJyDCzwlJ_JLB9qYs_PLYYrIDGce51y0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=onedice-bot

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/onedice

[Install]
WantedBy=multi-user.target
